# screen_brightness/ux581_brightness.path
[Unit]
Description=Watch /sys for Brightness Changes and Update Displays
ConditionPathExists=/proc/acpi/call
ConditionPathExists=/sys/class/backlight/intel_backlight/brightness
ConditionPathExists=/sys/class/backlight/intel_backlight/max_brightness
ConditionPathExists=/usr/local/sbin/ux581_brightness.sh
StartLimitIntervalSec=0

[Path]
PathModified=/sys/class/backlight/intel_backlight/brightness

[Install]
WantedBy=multi-user.target
